# download.cron
* 18 * * * python3 downloadExtractBhavcopy.py
 > tatus.log
